# contents-agent/scripts/sh.voicebox.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>sh.voicebox.watchdog</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/yoogeon/Agents/contents agent/scripts/voicebox-watchdog.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/yoogeon/Library/Logs/voicebox-watchdog.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/yoogeon/Library/Logs/voicebox-watchdog.log</string>
    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>
